# agent/s6-overlay/s6-rc.d/runtime-monitor/run
#!/command/execlineb -P
# Runtime Monitor Service - Continuously monitors traffic enforcement
# Runs as part of s6 longrun service

foreground {
    exec /home/agent/scripts/runtime-monitor.sh
}