# agent/s6-overlay/s6-rc.d/runtime-monitor/run
#!/usr/bin/execlineb -P
# Runtime Monitor Service - Continuously monitors traffic enforcement
exec /home/agent/scripts/runtime-monitor.sh