# agent/s6-overlay/s6-rc.d/startup-validation/up
#!/command/execlineb -P
# Run startup validation tests
foreground { /home/agent/scripts/startup-validation.sh }