# agent/s6-overlay/s6-rc.d/traffic-enforce/up
#!/command/execlineb -P
# Initialize iptables traffic enforcement rules
# Must run early in container startup before any user processes

foreground { /home/agent/scripts/traffic-enforcement.sh }
