# agent/s6-overlay/s6-rc.d/gpg-init/up
#!/command/execlineb -P
# S6 oneshot service to initialize GPG for commit signing
# Runs once at container startup after git-init

foreground {
  s6-setuidgid agent
  /home/agent/scripts/setup-gpg.sh
}
