# agent/s6-overlay/s6-rc.d/build-config/up
#!/bin/bash
/home/agent/scripts/setup-build-config.sh
