# agent/s6-overlay/s6-rc.d/git-init/up
#!/command/execlineb -P
# S6 oneshot service to initialize Git credentials
# Runs once at container startup before agent service

foreground { 
  /home/agent/scripts/setup-git-credentials.sh
}