# agent/s6-overlay/s6-rc.d/bash-init/up
#!/command/execlineb -P
# Setup bash configuration on container startup
# This configures .bashrc to auto-source environment variables

foreground { /home/agent/scripts/setup-bashrc.sh }