# agent/s6-overlay/s6-rc.d/metrics/run
#!/command/execlineb -P
python3 /home/agent/scripts/metrics_exporter.py
